# launchd/com.lazydog.creator-data-dashboard.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lazydog.creator-data-dashboard.collector</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lazydog/Library/Application Support/CreatorDataDashboard/scripts/mac_mini_update_and_sync.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/lazydog/Library/Application Support/CreatorDataDashboard</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>ThrottleInterval</key>
  <integer>60</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/lazydog/Library/Logs/CreatorDataDashboard/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lazydog/Library/Logs/CreatorDataDashboard/launchd.err.log</string>
</dict>
</plist>
